# rootfs/etc/supervisor.d/dbus.ini

[program:dbus]
priority=10
autostart=true
autorestart=true
directory=/
user=%(ENV_USER)s
command=bash -c "until [ -d \"${XDG_RUNTIME_DIR}\" ]; do sleep 0.5; done; dbus-daemon --system --nofork --nosyslog --nopidfile --address=\"${DBUS_SYSTEM_BUS_ADDRESS}\""
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s",XDG_RUNTIME_DIR="%(ENV_XDG_RUNTIME_DIR)s",DBUS_SYSTEM_BUS_ADDRESS="%(ENV_DBUS_SYSTEM_BUS_ADDRESS)s"
stopsignal=INT
stdout_logfile=/var/log/supervisor/dbus.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/dbus.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
