# rootfs/etc/supervisor.d/selkies-gstreamer.ini

[program:selkies-gstreamer]
autostart=true
autorestart=true
priority=20
user=root
command=/usr/bin/selkies-gstreamer.sh
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
stopsignal=INT
stdout_logfile=/var/log/supervisor/selkies-gstreamer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/selkies-gstreamer.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7